# roles/letsencrypt/templates/certbot.timer.j2
[Unit]
Description=Timer for certbot certificate renewal

[Timer]
OnBootSec=300
OnUnitActiveSec=1w

[Install]
WantedBy=timers.target
