# roles/backups/templates/backup.timer.j2
[Unit]
Description=Perform system backup

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
